# supervisord.conf
[supervisord]
nodaemon=true

[program:tribes2-server]
command=/bin/su - gameserv -c /home/gameserv/start-server
